# services/ssh/run
#!/bin/sh
exec /usr/sbin/sshd >>/var/log/sshd.log 2>&1
