# services/ssh/run
#!/usr/bin/env bash

ps cax | grep sshd > /dev/null
if [ $? -ne 0 ]; then
  /usr/sbin/sshd
fi
