# services/nginx/run
#!/bin/sh
exec /usr/sbin/nginx >>/var/log/nginx/nginx.log 2>&1
