# services/php-fpm/run
#!/bin/sh

exec php-fpm7.4 --nodaemonize -c /etc/php/7.4/fpm
