# services/php-fpm/run
#!/bin/sh

exec php-fpm7.0 --nodaemonize -c /etc/php/7.0/fpm
