# services/php-fpm/run
#!/bin/sh
exec php-fpm5.6 --nodaemonize -c /etc/php/5.6/fpm
